# .config/systemd/user/foot-server.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/foot --server=3 --login-shell
